# actividad2/crontab
*/5 * * * * bash /home/dave965/act2.sh
